# deploy/cronix.service
[Unit]
Description=Cronix Task Scheduler
Documentation=https://github.com/riverisagame/cronix
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/cronix
ExecStart=/opt/cronix/cronix serve -c /opt/cronix/config.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5

# 日志 → journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cronix

[Install]
WantedBy=multi-user.target
